# ubuntu/etc/systemd/system/docker.service.d/socks5-proxy.conf
[Service]
Environment="ALL_PROXY=socks5://proxy.ip:port"
Environment="NO_PROXY=localhost,127.0.0.1,www.artifactory.com,130.147.0.0/16,130.145.0.0/16"
